# assets/ffdyndns.service
[Unit]
Description=Simple dynamic dns service
After=network.target pdns.service

[Service]
ExecStart=/usr/bin/ffdyndns server
WorkingDirectory=/usr/lib/ffdyndns
Restart=on-abort

[Install]
WantedBy=multi-user.target
